# config/local/supervisord-qa2-1-tornado.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[group:tornadoes]
programs=tornado-8000

[program:tornado-8000]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8000
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info
