# config/local/supervisord-qa2.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[group:tornadoes]
programs=tornado-8000,tornado-8001,tornado-8002,tornado-8003,tornado-8004,tornado-8005,tornado-8006,tornado-8007,tornado-8008,tornado-8009,tornado-8010,tornado-8011,tornado-8012,tornado-8013,tornado-8014,tornado-8015,tornado-8016,tornado-8017,tornado-8018,tornado-8019

[program:tornado-8000]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8000
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8001]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8001
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8002]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8002
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8003]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8003
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8004]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8004
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8005]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8005
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8006]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8006
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8007]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8007
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8008]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8008
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8009]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8009
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8010]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8010
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8011]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8011
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8012]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8012
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8013]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8013
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8014]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8014
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8015]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8015
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8016]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8016
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8017]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8017
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8018]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8018
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info

[program:tornado-8019]
command=/opt/api_semantica/brainiak/virtualenv/bin/python -m brainiak.server --port=8019
directory=/var/www
autorestart=true
redirect_stderr=true    
stdout_logfile=/tmp/tornado-0.log
loglevel=info
